# arch-etc/profile.d/zdotdir.sh
export ZDOTDIR="$HOME/.config/zsh"
export GTK_THEME=Adwaita:dark
export EDITOR=vim
