# arch-etc/profile.d/zdotdir.sh
export ZDOTDIR="$HOME/.config/zsh"
export XINITRC="$HOME/.config/X11/xinitrc"
export GTK_THEME=Adwaita:dark
export EDITOR=vim
